# sys-apps/hide-hardware-info/files/hide-hardware-info.initd-r1
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="hide-hardware-info"
description="Hide hardware information to unprivileged users"
command="/usr/bin/hide-hardware-info"
